# roles/base-system/files/rkhunter-scan.service.j2
[Unit]
Description=Run rkhunter security scan
Documentation=man:rkhunter

[Service]
Type=oneshot
ExecStart=/usr/bin/rkhunter --checkall --skip-keypress --report-warnings-only
