# roles/base-system/templates/clamscan.service.j2
[Unit]
Description=ClamAV On-Demand Scanner Service

[Service]
Type=oneshot
ExecStart=/usr/bin/clamscan {{ base_system_clamscan_options }}

