# roles/base-system/templates/rkhunter-scan.timer.j2
[Unit]
Description=Run rkhunter scan periodically

[Timer]
OnCalendar={{ base_system_rkhunter_scan_interval }}
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
